# files/etc/systemd/system/nodepool.service
[Unit]
Description=Nodepool Service
After=syslog.target network.target

[Service]
Type=simple
User=nodepool
Group=nodepool
EnvironmentFile=-/etc/sysconfig/nodepool
ExecStart=/bin/sh -c "${PREFIX}/bin/nodepoold -d --no-builder ${OPTIONS}"

[Install]
WantedBy=multi-user.target
